# etc/supervisor/conf.d/my_test_app.conf
[program:my_test_app]
command=gunicorn -b localhost:8000 -w 4 wsgi:my_test_app
directory=/my_test_app/src
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
